# dashboard/dashboard.service
[Unit]
Description=My custom startup script

[Service]
ExecStart=/home/root/homelab/dashboard/start.sh

[Install]
WantedBy=multi-user.target
